# desktop.target
[Unit]
Description=Graphical desktop environment
Requires=graphical-session.target
